# systemd/system/arp-service.service
[Unit]
Description=Example systemd service.

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/arp-service.sh

[Install]
WantedBy=multi-user.target
