# network/if-up.d/static-arp
#! / Bin / sh
arp -i eth0 -s 193.191.177.254 ca:fe:c0:ff:ee:00
arp -i eth0 -s 193.191.177.219 a2:fc:bf:0c:33:8a
